# deploy/systemd/devdock.service
# ============================================================================
# DEPRECATED (Phase 999.2): DevDock now runs as a Docker Compose stack.
#
# This systemd unit is retained in-place (D-18) until the VPS cutover is
# confirmed successful. Once Docker Compose has been running cleanly in
# production for a verification window, this file (plus docker-compose.dev.yml
# removal and any remaining deprecated-in-place references) will be removed
# in a follow-up cleanup phase.
#
# Cutover runbook:
#   .planning/phases/999.2-devdock-self-containerization/999.2-HUMAN-UAT.md
#
# Cutover commands (summary):
#   sudo systemctl stop devdock
#   sudo systemctl disable devdock
#   # (update host nginx configs per deploy/nginx/*.conf changes)
#   sudo nginx -t && sudo systemctl reload nginx
#   cd /home/murx-dev/devdock && docker compose up -d --build
# ============================================================================

# deploy/systemd/devdock.service
# Install: sudo cp deploy/systemd/devdock.service /etc/systemd/system/
# Enable: sudo systemctl enable devdock
# Start: sudo systemctl start devdock
# Logs: journalctl -u devdock -f

[Unit]
Description=DevDock Development Platform
Documentation=https://github.com/your-repo/devdock
After=network.target postgresql.service docker.service
Requires=docker.service

[Service]
Type=simple
User=murx-dev
Group=murx-dev
WorkingDirectory=/home/murx-dev/devdock
ExecStart=/usr/bin/node .next/standalone/server.js
Restart=on-failure
RestartSec=5

# Environment
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
Environment=NODE_ENV=production
EnvironmentFile=/home/murx-dev/devdock/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/murx-dev/devdock/data
ReadWritePaths=/var/run/docker.sock
PrivateTmp=true
SupplementaryGroups=docker

# Resource limits for the DevDock process itself (D-17: ~512MB budget)
MemoryMax=512M
MemoryHigh=384M

[Install]
WantedBy=multi-user.target
